# deploy/systemd-hub.service
[Unit]
Description=Codex Bridge Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/codex-bridge
Environment=APP_ENV=prod
Environment=APP_HOST=127.0.0.1
Environment=APP_PORT=8088
ExecStart=/usr/local/bin/codex-bridge hub
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
